# systemd/vtech-trainer.service
[Unit]
Description=VTech Model Trainer Service
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=vtech
Group=vtech
WorkingDirectory=/opt/vtech
EnvironmentFile=/opt/vtech/.env
ExecStart=/opt/vtech/.venv/bin/python -m src.training.trainer --start rolling --end today
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vtech-trainer
MemoryMax=16G

[Install]
WantedBy=multi-user.target
